# roles/base/default_files/templates/etc/profile.d/10-environment.sh.j2
#
#  Global environment variables
#
export ENV='{{ subdomain_name }}'
