# roles/base/_base/templates/etc/profile.d/environment.sh.j2
#===============================================================================
#
#  E N V I R O N M E N T . S H
#
#  Global environment variables
#
#===============================================================================
export ENV='{{ environment_name }}'
